# supervisord.conf
[supervisord]
nodaemon=true

[program:script1]
command=python bot.py

[program:script2]
command=python main.py --camera 1

[program:script3]
command=python main.py --camera 2