# supervisord.conf
[supervisord]
nodaemon=true

[program:script1]
command=python /app/main.py --camera 1

[program:script2]
command=python /app/main.py --camera 2

[program:script3]
command=python /app/bot.py